# photo-frame-server.service
[Unit]
Description=Digital Photo Frame Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rpi-frame
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DISPLAY=:0
ExecStart=/usr/bin/node /home/pi/rpi-frame/dist/server/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

